# celeste-input-bridge.service
[Unit]
Description=Celeste input overlay bridge (evdev → WebSocket)
After=graphical-session.target

[Service]
ExecStart=/usr/bin/python3 %h/Games/celeste-overlay/celeste-input-bridge.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
